# docker/jetson/pserver.service.in
##regist service
#sed s/@CONFIG_PATH@/path_to_config/g pserver.pserver.in | sudo tee /etc/systemd/system/pserver.service
#sudo systemctl daemon-reload
#sudo systemctl start pserver.service
#sudo systemctl enable pserver.service

[Unit]
Description = pserver daemon

[Service]
ExecStart = /usr/bin/docker /usr/bin/pserver -c @CONFIG_PATH@
ExecStop = -/usr/bin/docker kill pserver
Restart = always
RestartSec = 3
Type = simple

[Install]
WantedBy = multi-user.target