# docker/jetson/pserver.service.in
##regist service
#sed -e "s/@CONFIG_PATH@/\/path\/to\/config.json/g" pserver.service.in | sudo tee /etc/systemd/system/pserver.service
#sudo systemctl daemon-reload
#sudo systemctl start pserver.service
#sudo systemctl enable pserver.service
#echo "xhost +" >> ~/.xsessionrc

[Unit]
Description = pserver daemon

[Service]
ExecStart = /usr/bin/pserver -c @CONFIG_PATH@
ExecStop = -/usr/bin/docker kill pserver

##docker container type
#ExecStart = /usr/bin/docker start -a pserver

Restart = always
RestartSec = 3
Type = simple

[Install]
WantedBy = multi-user.target
